# mushroom-client.service
[Unit]
Description=Mushroom Automation Client
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Desktop/raspClient2
Environment="PATH=/home/pi/Desktop/raspClient2/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/pi/Desktop/raspClient2/venv/bin/python3 boot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
